# deploy/db-backup.timer
[Unit]
Description=Daily PostgreSQL backup at 06:00 KST

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
